# config/switchberry/systemd/switchberry-apply-timing.service.d/raspap-switchberry.conf
[Service]
ExecStart=
ExecStart=/usr/local/sbin/raspap-switchberryctl _apply-timing
