# config/switchberry/systemd/switchberry-eth0-fixup.service.d/raspap-switchberry.conf
[Service]
ExecCondition=/usr/local/sbin/raspap-switchberryctl _require-tc-plane
ExecStart=
ExecStart=/usr/local/sbin/raspap-switchberryctl _eth0-fixup
