# config/switchberry/systemd/switchberry-sanity.service.d/raspap-switchberry.conf
[Service]
ExecStart=
ExecStart=/usr/local/sbin/raspap-switchberryctl _sanity
