# config/switchberry/systemd/switchberry-full-init.service.d/raspap-switchberry.conf
[Service]
ExecCondition=/usr/local/sbin/raspap-switchberryctl _require-tc-plane
